# infrastructure/playbooks/app/root/etc/services.d/nodejs/run
#!/bin/sh

exec node /var/www/app/index.js
