# infrastructure/playbooks/scheduler/root/etc/services.d/nomad/run
#!/bin/ash

exec s6-setuidgid nomad nomad agent -config /etc/nomad.d/server.hcl
